# server/zfsrecvsnap@.service
[Unit]
Description=ZFS snapshot receiver service

[Service]
ExecStart=/home/service/zfsrecvsnap/zfsrecvsnap.sh
StandardInput=socket
StandardOutput=socket
StandardError=journal
